# grafana/crontab
0 0 * * * /home/ec2-user/restartStats.sh
